# octavia-amphora-agent.service
[Unit]
Description=OpenStack Octavia Amphora Agent service
After=syslog.target network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/amphora-agent --config-file /etc/octavia/amphora-agent.conf
KillMode=process

[Install]
WantedBy=multi-user.target
